# dreamseed-api.service
[Unit]
Description=DreamSeed API Server (Gunicorn)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/home/won/projects/dreamseed_monorepo/dreamseed.env
WorkingDirectory=/home/won/projects/dreamseed_monorepo
User=won
Group=won

# UFW 규칙은 수동으로 설정됨 (포트 8000 허용됨)

# Gunicorn 실행
ExecStart=/usr/bin/env bash -lc '\
  cd "${WORKDIR}" && \
  source "${VENV}/bin/activate" && \
  exec gunicorn \
    --config "${WORKDIR}/gunicorn.conf.py" \
    --bind "0.0.0.0:${PORT}" \
    "${APP_MODULE}" \
'

# 무중단 재로딩(HUP)
ExecReload=/bin/kill -HUP $MAINPID

Restart=always
RestartSec=5
TimeoutStopSec=15
LimitNOFILE=65535

# 기본 보안 설정
NoNewPrivileges=yes
PrivateTmp=true

[Install]
WantedBy=multi-user.target
