# dreamseed-editor.service
[Unit]
Description=DreamSeed Editor (Express proxy + static)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Set WORKDIR in /etc/default/dreamseed-editor, e.g., WORKDIR=/opt/dreamseed
EnvironmentFile=/etc/default/dreamseed-editor
WorkingDirectory=${WORKDIR}
ExecStart=/usr/bin/env node proxy_server.js
Restart=on-failure
RestartSec=2
Environment=NODE_ENV=production

[Install]
WantedBy=default.target
