# alert-threader/ops-services-alert-threader-go-envfile.service
[Unit]
Description=DreamSeed Alert Threader - Go (Blocks, Persistent, EnvironmentFile, go-redis)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/alert-threader-go-redis
EnvironmentFile=/etc/alert-threader.env

# 빌드(최초/업데이트 시)
ExecStartPre=/usr/bin/env bash -lc '\
  mkdir -p /var/lib/alert-threader && chown www-data:www-data /var/lib/alert-threader && \
  cd /opt/alert-threader-go-redis && \
  [ -f go.mod ] || go mod init alert-threader-go-redis && \
  go mod tidy && \
  go build -o threader .'

User=www-data
Group=www-data
ExecStart=/usr/bin/env bash -lc '/opt/alert-threader-go-redis/threader'
Restart=always
RestartSec=5
TimeoutStopSec=15

# --- Hardening ---
NoNewPrivileges=yes
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ProtectControlGroups=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectClock=yes
LockPersonality=yes
CapabilityBoundingSet=~CAP_SYS_ADMIN CAP_NET_ADMIN CAP_SYS_MODULE CAP_SYS_PTRACE CAP_SETUID CAP_SETGID
RestrictSUIDSGID=yes
RestrictRealtime=yes
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native
SystemCallFilter=@system-service @basic-io @network-io ~@mount ~@privileged

# 로그 설정
StandardOutput=journal
StandardError=journal
SyslogIdentifier=alert-threader-go

[Install]
WantedBy=multi-user.target

