# alert-threader/ops/services/guard-brief.service
[Unit]
Description=Guard-Window Slack Briefing Service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/guard-brief.env
ExecStart=/usr/local/sbin/guard_window_brief.sh
User=root
Group=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target


