# network_automation/scripts/update-cloudflare-real-ip.timer
[Unit]
Description=Daily Cloudflare CIDR sync for nginx real_ip

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=600

[Install]
WantedBy=timers.target
