# alert-threader/systemd/alert-threader-python.service
[Unit]
Description=DreamSeed Alert Threader (Python)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/alert-threader/python

# 환경 변수
Environment=SLACK_BOT_TOKEN=xoxb-xxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxxx
Environment=SLACK_CHANNEL=C0123456789
Environment=ENVIRONMENT=production
Environment=BIND_HOST=0.0.0.0
Environment=BIND_PORT=9009

# 실행 명령
ExecStart=/usr/bin/env bash -lc 'cd /opt/alert-threader/python && python3 -m pip install -r requirements.txt && python3 app.py'

Restart=always
RestartSec=5
TimeoutStopSec=15

# 보안 설정
NoNewPrivileges=yes
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

# 로그 설정
StandardOutput=journal
StandardError=journal
SyslogIdentifier=alert-threader

[Install]
WantedBy=multi-user.target

